# 2config/bspwm/scripts.old/99-reload-monitor.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/bin/su anoop --command='systemctl --user start bspwm-reload.service'"
 
